# roles/euresys-framegrabber/files/euresys_kernel_modules.service
[Unit]
After=ssh.service

[Service]
ExecStart=/bin/bash /usr/local/bin/install_euresys_modules.sh

[Install]
WantedBy=default.target
